# Bin/systemd/rml_seman.service
[Unit]
Description=RML seman daemon
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
User=sokirko
Restart=on-failure
RestartSec=5s
Environment="RML=/home/sokirko/RML"
ExecStart=/home/sokirko/RML/Bin/SemanDaemon --host 63.250.59.227 --port 8081


[Install]
WantedBy=multi-user.target
